# schedulers/com.samc.hyperdeck-archiver.prune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Daily NAS prune (delete folders older than retention.days). Install:
     cp to ~/Library/LaunchAgents/ then
     launchctl load -w ~/Library/LaunchAgents/com.samc.hyperdeck-archiver.prune.plist
     Replace @@VENV_PYTHON@@ and @@INSTALL_DIR@@ with absolute paths. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.samc.hyperdeck-archiver.prune</string>
    <key>ProgramArguments</key>
    <array>
        <string>@@VENV_PYTHON@@</string>
        <string>@@INSTALL_DIR@@/run.py</string>
        <string>--config</string>
        <string>@@INSTALL_DIR@@/config.yaml</string>
        <string>prune</string>
    </array>
    <key>WorkingDirectory</key>
    <string>@@INSTALL_DIR@@</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>@@INSTALL_DIR@@/logs/prune.out.log</string>
    <key>StandardErrorPath</key>
    <string>@@INSTALL_DIR@@/logs/prune.err.log</string>
</dict>
</plist>
